# bads-camera.service
[Unit]
Description=BADS Flask Camera Server for ESP32
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bads-rpi
ExecStart=/usr/bin/python3 /opt/bads-rpi/http/flask_camera_server.py
Restart=always
RestartSec=10

# Environment variables
Environment=PYTHONPATH=/opt/bads-rpi

[Install]
WantedBy=default.target